# cron.d/killVpn
# Disconnect VPN at 18:00
SHELL=/bin/bash
00   18 *   *   *   root /usr/bin/killall openconnect
